# .local/share/systemd/user/tilda.service
[Unit]
Description=Tilda console

[Service]
Type=simple
Restart=always
RestartSec=3
TimeoutStartSec=0
ExecStart=/usr/bin/tilda

[Install]
WantedBy=default.target
